# ft8tx.service
[Unit]
Description=FT8 Transmitter Service

[Service]
ExecStart=/usr/local/bin/ft8
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target